# infra/systemd/poi-backup.service
[Unit]
Description=POI stack MariaDB backup job

[Service]
Type=oneshot
WorkingDirectory=__PROJECT_ROOT__
ExecStart=__PROJECT_ROOT__/scripts/backup.sh
