# infra/systemd/poi-backup.service
[Unit]
Description=POI stack MariaDB backup job

[Service]
Type=oneshot
WorkingDirectory=%h/0_opencode/application_poi
ExecStart=%h/0_opencode/application_poi/scripts/backup.sh
