# infra/quadlet/poi-api.container
[Unit]
Description=POI API service
After=poi-db.service
Requires=poi-db.service

[Container]
ContainerName=poi-api
Image=localhost/poi-api:latest
Network=poi.network
EnvironmentFile=%h/.config/poi-stack/poi.env
PublishPort=127.0.0.1:3001:3001

[Service]
Restart=always
TimeoutStartSec=120

[Install]
WantedBy=default.target
